# aios-theme.service
[Unit]
Description=AIOS Theming Agent
Before=gdm.service

[Service]
ExecStart=/usr/bin/install-default-theme

[Install]
WantedBy=default.target

